# .config/systemd/user/win-to-clipboard.service
[Unit]
Description=Sync windows clipboard to wayland

[Service]
ExecStart=%h/.local/bin/win-to-clipboard
